# AppServices/SD.IdentitySystem.AppService.Host(gRPC)/Installers/SD.IdentitySystem.AppService(gRPC).service
[Unit]
Description=SD.IdentitySystem.AppService(gRPC)

[Service]
ExecStart=/usr/bin/dotnet /home/lee/Publish/身份认证系统服务(gRPC)/SD.IdentitySystem.AppService.Host.dll

[Install]
WantedBy=multi-user.target
